# apps/api/resources/deploy/comecyt-queue.conf
[program:comecyt-queue]
; Supervisor config para el worker de colas (emails, notificaciones async)
; Copiar a: /etc/supervisor/conf.d/comecyt-queue.conf
; Luego: supervisorctl reread && supervisorctl update && supervisorctl start comecyt-queue

command=php /var/www/comecyt/apps/api/artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www/comecyt/apps/api
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/comecyt/apps/api/storage/logs/queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopasgroup=true
killasgroup=true
