# apps/api/resources/deploy/comecyt-reverb.conf
[program:comecyt-reverb]
; Supervisor config para Laravel Reverb (WebSockets de notificaciones en tiempo real)
; Copiar a: /etc/supervisor/conf.d/comecyt-reverb.conf
; Luego: supervisorctl reread && supervisorctl update && supervisorctl start comecyt-reverb
;
; Reverb escucha en REVERB_PORT (default 8080). En producción, expón este puerto
; detrás de Nginx con proxy_pass + Upgrade:websocket headers, o usa una URL pública
; tipo wss://ws.dominio.edomex.gob.mx → 127.0.0.1:8080.

command=php /var/www/comecyt/apps/api/artisan reverb:start --host=0.0.0.0 --port=8080
directory=/var/www/comecyt/apps/api
user=www-data
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/comecyt/apps/api/storage/logs/reverb.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
stopasgroup=true
killasgroup=true
