# apps/api/resources/deploy/comecyt-scheduler.conf
[program:comecyt-scheduler]
; Supervisor config para el scheduler de Laravel (convocatorias, notificaciones, etc.)
; Copiar a: /etc/supervisor/conf.d/comecyt-scheduler.conf
; Luego: supervisorctl reread && supervisorctl update && supervisorctl start comecyt-scheduler

command=php /var/www/comecyt/apps/api/artisan schedule:work
directory=/var/www/comecyt/apps/api
user=www-data
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/comecyt/apps/api/storage/logs/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopasgroup=true
killasgroup=true
