# docs/deploy-templates/supervisor-comecyt-next.conf
# =============================================================================
# Supervisor — Next.js Frontend Server
# =============================================================================
# Ubicación: /etc/supervisor/conf.d/comecyt-next.conf
# =============================================================================

[program:comecyt-next]
command=/usr/bin/npm run start
directory=/var/www/comecyt/apps/web
user=www-data
group=www-data
autostart=true
autorestart=true
startretries=5
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/comecyt-next.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=NODE_ENV="production",PORT="3000"
