# docs/deploy-templates/supervisor-comecyt-reverb.conf
# =============================================================================
# Supervisor — Reverb WebSocket Server
# =============================================================================
# Ubicación: /etc/supervisor/conf.d/comecyt-reverb.conf
# Activar:   sudo supervisorctl reread && sudo supervisorctl update
# =============================================================================

[program:comecyt-reverb]
command=/usr/bin/php /var/www/comecyt/apps/api/artisan reverb:start --host=127.0.0.1 --port=8080
directory=/var/www/comecyt/apps/api
user=www-data
group=www-data
autostart=true
autorestart=true
startretries=5
stopwaitsecs=10
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/var/log/comecyt-reverb.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=APP_ENV="production"
