# docs/deploy-templates/supervisor-comecyt-queue.conf
# =============================================================================
# Supervisor — Queue Workers (procesamiento asíncrono)
# =============================================================================
# Procesa: envío de emails, notificaciones push, generación de PDFs pesados
# Ubicación: /etc/supervisor/conf.d/comecyt-queue.conf
# =============================================================================

[program:comecyt-queue]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/comecyt/apps/api/artisan queue:work --sleep=3 --tries=3 --max-time=3600 --backoff=10
directory=/var/www/comecyt/apps/api
user=www-data
group=www-data
numprocs=2
autostart=true
autorestart=true
startretries=5
stopwaitsecs=3600
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/var/log/comecyt-queue.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=APP_ENV="production"
